# templates/default/upstart/jently.conf.erb
start on runlevel [2345]
stop on runlevel [06]

chdir <%= @server_dir %>
setuid jently

respawn
exec bundle exec ruby jently.rb
